# contrib/callme.conf
[program:callme]
command=/usr/bin/php callme.php
directory=/var/www/html/callme
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/www/html/callme/logs/daemon.log